# Linux/auditd.yaml
## First rule - delete all
- -D

#Contents:

#Configuration
#Interesting process activity
#Suspicious activity
#Network activity
#Process activity

#--------------------------------------------------------------------
#                       Configuration
#--------------------------------------------------------------------

## Things that could affect time
- -a always,exit -F arch=b32 -S adjtimex,settimeofday,stime -F key=time-change
- -a always,exit -F arch=b64 -S adjtimex,settimeofday -F key=time-change
- -a always,exit -F arch=b32 -S clock_settime -F a0=0x0 -F key=time-change
- -a always,exit -F arch=b64 -S clock_settime -F a0=0x0 -F key=time-change
- -w /etc/localtime -p wa -k time-change

## Things that could affect system locale 6
- -a always,exit -F arch=b32 -S sethostname,setdomainname -F key=system-locale
- -a always,exit -F arch=b64 -S sethostname,setdomainname -F key=system-locale
- -w /etc/issue -p wa -k system-locale
- -w /etc/issue.net -p wa -k system-locale
- -w /etc/hosts -p wa -k system-locale
- -w /etc/hostname -p wa -k system-locale

#- -w /etc/sysconfig/network -p wa -k system-locale 12
- -a always,exit -F dir=/etc/NetworkManager/ -F perm=wa -F key=system-locale

## Modprobe configuration
- -w /etc/modprobe.conf -p wa -k modprobe

## watch for containers that may change their configuration 14
- -a always,exit -F arch=b32 -S unshare,setns -F key=container-config
- -a always,exit -F arch=b64 -S unshare,setns -F key=container-config

# Configure if using SElinux
#- -w /etc/selinux -p wa -k selinux-changes

#Configure with the location of auditd or go-audit dirs
#- -w /var/log/audit/ -k audit-logs

## User file changes 16
- -w /etc/group -p wa -k auth-changes
- -w /etc/passwd -p wa -k auth-changes
- -w /etc/gshadow -p wa -k auth-changes
- -w /etc/shadow -p wa -k auth-changes
- -w /etc/security/opasswd -p wa -k auth-changes

## Sudoers file changes 21
- -w /etc/sudoers -p wa -k auth-changes

## SSH service config changes 22

- -w /etc/ssh/sshd_config -k sshd

## Cron configuration & scheduled jobs 23
- -w /etc/cron.allow -p wa -k cron
- -w /etc/cron.deny -p wa -k cron
- -w /etc/cron.d/ -p wa -k cron
- -w /etc/cron.daily/ -p wa -k cron
- -w /etc/cron.hourly/ -p wa -k cron
- -w /etc/cron.monthly/ -p wa -k cron
- -w /etc/cron.weekly/ -p wa -k cron
- -w /etc/crontab -p wa -k cron
- -w /var/spool/cron/crontabs/ -k cron

## unsuccessful creation 32
- -a always,exit -F arch=b32 -S creat,link,mknod,mkdir,symlink,mknodat,linkat,symlinkat -F exit=-EACCES -F key=creation
- -a always,exit -F arch=b64 -S mkdir,creat,link,symlink,mknod,mknodat,linkat,symlinkat -F exit=-EACCES -F key=creation
- -a always,exit -F arch=b32 -S link,mkdir,symlink,mkdirat -F exit=-EPERM -F key=creation
- -a always,exit -F arch=b64 -S mkdir,link,symlink,mkdirat -F exit=-EPERM -F key=creation

## unsuccessful open 36
- -a always,exit -F arch=b32 -S open,openat,open_by_handle_at -F exit=-EACCES -F key=open
- -a always,exit -F arch=b64 -S open,openat,open_by_handle_at -F exit=-EACCES -F key=open
- -a always,exit -F arch=b32 -S open,openat,open_by_handle_at -F exit=-EPERM -F key=open
- -a always,exit -F arch=b64 -S open,openat,open_by_handle_at -F exit=-EPERM -F key=open

## unsuccessful close 40
- -a always,exit -F arch=b32 -S close -F exit=-EIO -F key=close
- -a always,exit -F arch=b64 -S close -F exit=-EIO -F key=close

## unsuccessful modifications 42
- -a always,exit -F arch=b32 -S rename -S renameat -S truncate -S chmod -S setxattr -S lsetxattr -S removexattr -S lremovexattr -F exit=-EACCES -F key=mods
- -a always,exit -F arch=b64 -S rename -S renameat -S truncate -S chmod -S setxattr -S lsetxattr -S removexattr -S lremovexattr -F exit=-EACCES -F key=mods
- -a always,exit -F arch=b32 -S rename -S renameat -S truncate -S chmod -S setxattr -S lsetxattr -S removexattr -S lremovexattr -F exit=-EPERM -F key=mods
- -a always,exit -F arch=b64 -S rename -S renameat -S truncate -S chmod -S setxattr -S lsetxattr -S removexattr -S lremovexattr -F exit=-EPERM -F key=mods

## unsuccessful deletion 46
- -a always,exit -F arch=b32 -S unlink,rmdir,unlinkat -F exit=-EACCES -F key=delete
- -a always,exit -F arch=b64 -S rmdir,unlink,unlinkat -F exit=-EACCES -F key=delete
- -a always,exit -F arch=b32 -S unlink,rmdir,unlinkat -F exit=-EPERM -F key=delete
- -a always,exit -F arch=b64 -S rmdir,unlink,unlinkat -F exit=-EPERM -F key=delete

#--------------------------------------------------------------------
#                       Interesting process activity 51
#--------------------------------------------------------------------

- -w /bin/bash -p x -k interesting-bin
- -w /bin/mknod -p x -k interesting-bin
- -w /bin/nc -p x -k interesting-bin
- -w /bin/nc.openbsd -p x -k interesting-bin
- -w /bin/netcat -p x -k interesting-bin
- -w /bin/systemctl -p x -k systemd
- -w /usr/bin/base64 -p x -k interesting-bin
- -w /usr/bin/curl -p x -k interesting-bin
- -w /usr/bin/ftp -p x -k interesting-bin
- -w /usr/bin/gcc -p x -k interesting-bin
- -w /usr/bin/ncat -p x -k interesting-bin
- -w /usr/bin/nmap -p x -k interesting-bin
- -w /usr/bin/rawshark -p x -k interesting-bin
- -w /usr/bin/rdesktop -p x -k interesting-sbin
- -w /usr/bin/scp -p x -k interesting-bin
- -w /usr/bin/socat -p x -k interesting-bin
- -w /usr/bin/ssh -p x -k interesting-bin
- -w /usr/bin/strace -p x -k interesting-bin
- -w /usr/bin/tcpdump -p x -k interesting-bin
- -w /usr/bin/wget -p x -k interesting-bin
- -w /usr/bin/whois -p x -k interesting-bin

# Sbin interesting activity 74

- -w /sbin/auditctl -p x -k audittools
- -w /sbin/auditd -p x -k audittools
- -w /sbin/ifconfig -p x -k interesting-sbin
- -w /sbin/iptables -p x -k interesting-sbin
- -w /usr/sbin/addgroup -p x -k interesting-sbin
- -w /usr/sbin/adduser -p x -k interesting-sbin
- -w /usr/sbin/groupadd -p x -k interesting-sbin
- -w /usr/sbin/groupmod -p x -k interesting-sbin
- -w /usr/sbin/tcpdump -p x -k interesting-sbin
- -w /usr/sbin/traceroute -p x -k interesting-sbin
- -w /usr/sbin/useradd -p x -k interesting-sbin
- -w /usr/sbin/usermod -p x -k interesting-sbin

#--------------------------------------------------------------------
#                       Suspicious activity
#--------------------------------------------------------------------

# Privliege elevation 86

- -w /bin/su -p x -k elevation
- -w /usr/bin/sudo -p x -k elevation

## An admin may be rooting around
- -a always,exit -F dir=/home -F uid=0 -F auid>=1000 -F auid!=4294967295 -C auid!=obj_uid -F key=root

## Ptrace code injection - could be suspicious or just debugging

#-a always,exit -F arch=b32 -S ptrace -F key=tracing (55)
- -a always,exit -F arch=b32 -S ptrace -F a0=0x4 -F key=code-injection
- -a always,exit -F arch=b64 -S ptrace -F a0=0x4 -F key=code-injection
- -a always,exit -F arch=b32 -S ptrace -F a0=0x5 -F key=data-injection
- -a always,exit -F arch=b64 -S ptrace -F a0=0x5 -F key=data-injection
- -a always,exit -F arch=b32 -S ptrace -F a0=0x6 -F key=register-injection
- -a always,exit -F arch=b64 -S ptrace -F a0=0x6 -F key=register-injection
- -a always,exit -F arch=b64 -S ptrace -F key=tracing

## These rules watch for kernel module insertion
- -w /sbin/kmod -p x -k modules
- -w /sbin/insmod -p x -k modules
- -w /sbin/rmmod -p x -k modules
- -w /sbin/modprobe -p x -k modules
- -a always,exit -F arch=b32 -S init_module,finit_module -F key=module-load
- -a always,exit -F arch=b64 -S init_module,finit_module -F key=module-load
- -a always,exit -F arch=b32 -S delete_module -F key=module-unload
- -a always,exit -F arch=b64 -S delete_module -F key=module-unload

## These are rules are to locate poorly written programs.
## Its never planned to waste time on a syscall with incorrect parameters
## This is more of a debugging step than something people should run with
## in production. (70-71)
- -a never,exit -F arch=b64 -S rt_sigreturn
- -a always,exit -S all -F exit=-EINVAL -F key=einval-retcode

#--------------------------------------------------------------------
#                       Network activity
#--------------------------------------------------------------------

## All network syscalls making or recieving network connections (51)

- -a always,exit -F arch=b64 -S accept -F key=network
- -a always,exit -F arch=b64 -S bind -F key=network
- -a always,exit -F arch=b64 -S connect -F key=network
- -a always,exit -F arch=b64 -S listen -F key=network

#--------------------------------------------------------------------
#                       Process activity
#--------------------------------------------------------------------

## If you are on a 64 bit platform, everything _should_ be running
## in 64 bit mode. This rule will detect any use of the 32 bit syscalls
## because this might be a sign of someone exploiting a hole in the 32
## bit API.

- -a always,exit -F arch=b32 -S all -F key=32bit-API

# All process and command execution
- -a always,exit -F arch=b64 -S execve -F key=process
- -a always,exit -F arch=b64 -S execveat -F key=process

#Lastly
#- -e 1
